# systemd/oscmix.service
[Unit]
Description=oscmix backend for RME Fireface interfaces
Documentation=https://github.com/michaelforney/oscmix
After=sound.target
# Started on hotplug via the udev rule (SYSTEMD_USER_WANTS); stops again
# when the pull-in from the device disappears.
StopWhenUnneeded=yes
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
# oscmix-session sends READY=1 once the backend is up and the routing
# has been applied; verification then runs in the background.
Type=notify
ExecStart=%h/.local/bin/oscmix-session
Restart=on-failure
RestartSec=3
# Exit code 2 means a routing.conf error -- restarting cannot fix that.
RestartPreventExitStatus=2
# Device wait (30s) + port wait (10s) + verification, with headroom.
TimeoutStartSec=75
TimeoutStopSec=10
NoNewPrivileges=yes

[Install]
WantedBy=default.target
